# systemd/aar-clip-server.service
[Unit]
Description=AAR clip server — HTTP file server for finished archive clips
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/aar-clip-server
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
