# systemd/aar-sink-stream-music.service
[Unit]
Description=AAR sink stream (music) — null-sink → ffmpeg → HTTP MP3 fan-out
After=pipewire-pulse.socket
Wants=pipewire-pulse.socket

[Service]
Type=simple
ExecStart=%h/.local/bin/aar-sink-stream --sink aar-music --port 7772
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
