# systemd/agent-audio-relay.service
[Unit]
Description=agent-audio-relay — deliver agent TTS audio to playback target
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/agent-audio-relay
# Backend: ssh-termux (default) or mpv
Environment=RELAY_BACKEND=ssh-termux
# ssh-termux backend settings:
Environment=RELAY_SSH_HOST=p8ar
# mpv backend settings (uncomment to use):
# Environment=RELAY_BACKEND=mpv
# Environment=RELAY_MPV_SOCKET=/tmp/mpv-relay.sock
Environment=RELAY_WATCH_DIRS=/tmp/openclaw:/tmp
# Per-user queue dir so multiple OS users on one host don't collide on /tmp.
Environment=RELAY_QUEUE_DIR=%t/agent-audio-relay-queue
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
