# systemd/opencode-tts-watcher.service
[Unit]
Description=agent-audio-relay — OpenCode TTS watcher
After=network.target
Wants=agent-audio-relay.service

[Service]
Type=simple
# systemd's user PATH is minimal — if `opencode` or `edge-tts` live in
# ~/.nvm or ~/.local/bin, pin them explicitly:
# Environment=RELAY_OPENCODE_BIN=%h/.nvm/versions/node/<ver>/bin/opencode
# Environment=RELAY_EDGE_TTS_BIN=%h/.local/bin/edge-tts
#
# To use the OpenAI engine instead of Edge:
# Environment=OPENCODE_TTS_ENGINE=openai
# Environment=OPENCODE_TTS_VOICE=marin
# Environment=OPENCODE_TTS_OPENAI_PYTHON=%h/.pyenv/versions/<env>/bin/python3
# OPENAI_API_KEY must be in the unit env. One pattern: source a shell secrets file:
# ExecStart=/bin/bash -c 'set -a; . %h/.config/shell/secrets.zsh; exec %h/.local/bin/opencode-tts-hook'

ExecStart=%h/.local/bin/opencode-tts-hook
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
